# deploy/digitalocean/boxhaven-backend-backup.service
[Unit]
Description=Back up boxhaven backend state
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
EnvironmentFile=-/opt/boxhaven/app/deploy/digitalocean/.env.production
ExecStart=/usr/local/bin/boxhaven-backend-backup
